# Provision/files/supervisor/conf.d/example-daemon.conf
[program:unique-iden]
command=/usr/bin/php /vagrant/web/application/bin/cli.php -a module/controller/action
autorestart=true
autostart=true
user=www-data
startretries=100000000